# tests/integration/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:postgresql]
command=/usr/lib/postgresql/17/bin/postgres -D /var/lib/postgresql/data -c listen_addresses=localhost
user=postgres
priority=10
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log

[program:ocrmypdf]
command=gunicorn --bind 0.0.0.0:5000 --workers 2 --timeout 300 --chdir /app/ocrmypdf api:app
priority=20
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log

[program:stt]
command=stt serve --host 0.0.0.0 --port 8001
priority=20
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log

[program:mrdocument-service]
command=python -m mrdocument.server
directory=/app/service
priority=30
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log
environment=PYTHONPATH="/app/service"

[program:watcher-v2]
command=python3 app.py
directory=/app/watcher
priority=40
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log
